# overlay/etc/udev/rules.d/10-rfkill.rules
SUBSYSTEM=="rfkill", ATTR{name}=="ap6255", RUN+="/usr/local/bin/toggle-blue"
